# sh/fvq
#!/bin/bash
#
# chkconfig: 2345 55 25
# description: fvq service.
fvq_home=/home/ec2-user/FVQ
logs=$fvq_home/logs

########################################
# start fvq service
########################################
function start_service() {
  echo "Start fvq service..."
  nohup $fvq_home/sh/fvq-polling.sh 2> $logs/error.log &
}

########################################
# stop fvq service
########################################
function stop_service() {
  echo "Wait for fvq service stop..."
  pid=`ps aux | grep fvq-polling.sh |grep -v grep | awk '{print $2}'`
  if [ ! -z $pid ]; then
    kill $pid
  fi
}

########################################
# restart fvq service
########################################
function restart_service() {
  stop_service
  start_service
}


# main process
########################################
case "$1" in
  start)
    start_service
    ;;
  restart)
    restart_service
    ;;
  stop)
    stop_service
    ;;
  *)
    echo $"Usage: fvq {start|stop}"
esac
